# services/chat-intake/scripts/chat-intake.service
# DO NOT `systemctl restart chat-intake.service` directly after editing code.
# The runtime at /opt/forgingapps/chat-intake/ is a separate copy of the repo.
# Always run `./scripts/deploy-chat-intake.sh` from the repo root, which rsyncs
# the source, restarts both this service and the sweeper timer, and runs a
# drift check.
[Unit]
Description=ForgingApps chat-intake FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/forgingapps/chat-intake
User=alpharius
Group=alpharius
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/alpharius/.hermes
ExecStart=/opt/forgingapps/chat-intake/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/forgingapps/chat-intake.env

[Install]
WantedBy=multi-user.target
